# images/asteroid/files/system_files/etc/profile.d/99-aliases.sh
#!/usr/bin/env bash

alias dbox="distrobox"
alias ls="eza --icons"
alias cat="bat"