# images/nova/files/system_files/etc/profile.d/99-aliases.sh
#!/usr/bin/env bash

alias ls='eza -lah --group-directories-first'
alias cat='bat'